# servidorbackups/crontab



00  8   * * *   root  bash  /backupsscripts/servidorbackups/rsyncGodaddy.sh


